# crontab/rainbow
# rainbow images

@reboot pi /usr/local/bin/python3.6 -m karmapi.rainbow /home/pi/karmapi/rainbow


