# crontab/sense
@reboot pi /usr/bin/python3 -m karmapi.sense --path /home/pi/karmapi/sensehat

# crashes at 1 second past midnight, so start a new day
0 0 * * * pi /usr/bin/python3 -m karmapi.sense --path /home/pi/karmapi/sensehat
