# crontab/sense
@reboot pi /usr/local/bin/python3.6 -m karmapi.sense --path /home/pi/karmapi/sensehat --sleep 5
