# crontab/tankrain
DISPLAY=:0
PIG=hat
@reboot pi /usr/local/bin/python3.6 -m karmapi.tankrain /home/pi/karmapi/camera
